# lib/tasks/unicorn/unicorn.erb
#!/bin/sh
#
# chkconfig: - 85 15
# description: unicorn

. /etc/rc.d/init.d/functions

RAILS_USER=<%= ENV['USER'] %>
RAILS_ENV=<%= ENV['RAILS_ENV'] || 'development' %>
RAILS_ROOT=<%= ENV['RAILS_ROOT'] || Rails.root %>

PID=${RAILS_ROOT}/tmp/pids/unicorn.pid
CONFIG=${RAILS_ROOT}/config/unicorn.rb

export PATH="${PATH}":/usr/local/bin

start() {
  daemon --user=${RAILS_USER} --pidfile=$PID bundle exec unicorn -c ${CONFIG} -E ${RAILS_ENV} -D
}

stop() {
  killproc -p ${PID} -QUIT
}

cd ${RAILS_ROOT} || exit 1

case $1 in

  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo >&2 "Usage: $0 <start|stop|restart>"
    exit 1
    ;;
esac
