# itamae/templates/etc/init.d/god.erb
#!/bin/bash
#
# god       Startup script for God monitoring tool.
#
# chkconfig: - 85 15
# description: god monitors your system
#              

PID=/var/run/god.pid
CONFIG=/etc/god/master.conf
LOG=/var/log/god.log
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
RETVAL=0

case "$1" in
    start)
      echo 'Starting god'
      god -P $PID -c $CONFIG -l $LOG --log-level info
      RETVAL=$?
  ;;
    stop)
      if test -e $PID ; then
        god quit
      fi
      RETVAL=$?
  ;;
    restart)
      god quit
      echo 'Starting god'
      god -P $PID -c $CONFIG -l $LOG --log-level info
      RETVAL=$?
  ;;
    status)
      god status
      RETVAL=$?
  ;;
    *)
      echo "Usage: god {start|stop|restart|status}"
      exit 1
  ;;
esac

exit $RETVAL
